# crontab
0 10,12,14,16,18,20 * * * python3 /root/app.py